# crontab
*/1 * * * * cd /usr/bin/python   /app/script.py >> cron_log.log 2>&1